# dansal.service
[Unit]
Description=Dansal Calendar Server
After=network-online.target

[Service]
Type=simple
User=dansal
Group=dansal
ExecStartPre=/usr/lib/dansal/preflight
ExecStart=/usr/bin/dansal --config /etc/dansal/config.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5
TimeoutStopSec=35

# systemd creates /run/dansal (socket), /var/lib/dansal (database + images)
# and /etc/dansal (config) with correct ownership before starting the service
RuntimeDirectory=dansal
RuntimeDirectoryMode=0750
StateDirectory=dansal
StateDirectoryMode=0750
ConfigurationDirectory=dansal
ConfigurationDirectoryMode=0750

# Sandboxing
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
LockPersonality=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallFilter=@system-service
SystemCallArchitectures=native
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
