# dansal.service
[Unit]
Description=Dansal Calendar Server
After=network.target

[Service]
Type=simple
User=dansal
Group=dansal
ExecStart=/usr/bin/dansal --config /etc/dansal/config.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
TimeoutStopSec=35

# systemd creates /run/dansal (socket), /var/lib/dansal (database + images)
# and /etc/dansal (config) with correct ownership before starting the service
RuntimeDirectory=dansal
RuntimeDirectoryMode=0750
StateDirectory=dansal
StateDirectoryMode=0750
ConfigurationDirectory=dansal
ConfigurationDirectoryMode=0750

[Install]
WantedBy=multi-user.target
