# packaging/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Ensure state directories exist with correct ownership.
        # (systemd StateDirectory= handles this at runtime, but we need
        #  the images subdir too which systemd doesn't create.)
        install -d -m 750 -o dansal -g dansal /var/lib/dansal
        install -d -m 750 -o dansal -g dansal /var/lib/dansal/images

        # Ensure the database file is owned by the service user even if it was
        # previously created as root (e.g. during testing).
        touch /var/lib/dansal/calendar.db
        chown dansal:dansal /var/lib/dansal/calendar.db
        chmod 640           /var/lib/dansal/calendar.db

        # Secure the config file installed by dpkg.
        chown root:dansal /etc/dansal/config.yaml
        chmod 640         /etc/dansal/config.yaml

        install -d -m 750 -o dansal -g dansal /var/lib/dansal-web

        systemctl daemon-reload

        if [ -z "$2" ]; then
            # Fresh install — enable and start API; web requires manual config first.
            systemctl enable dansal.service
            systemctl start  dansal.service || true
            systemctl enable dansal-web.service
            echo "dansal-web installed. Edit /etc/dansal/web.yaml (set domain + dansal_url), then: systemctl start dansal-web"
        else
            # Upgrade — restart both services.
            systemctl try-restart dansal.service     || true
            systemctl try-restart dansal-web.service || true
        fi
        ;;
esac
